# 96luks-keyfile/luks-keyfile.target
[Unit]
Description=KEYFILE with LUKS
Before=sysinit.target dracut-pre-mount.service
After=systemd-ask-password-console.path
